# scripts/ryzen_adj_power_profiles/v3_power_profiles.service
[Unit]
Description="Automatically set TDP at different power levels."

[Service]
Type=simple
StandardOutput=journal
ExecStart=/path/to/power_adjust.sh

[Install]
WantedBy=default.target
